# docker-files/etc/supervisor/conf.d/supervisor-nginx.conf
[program:nginx]
command=/usr/local/openresty/nginx/sbin/nginx
priority=201
stdout_logfile=syslog
stderr_logfile=syslog