# docker-files/etc/supervisor/conf.d/supervisor-node.conf
[program:node]
command=/usr/local/bin/node /usr/src/app/app.js
priority=200
stdout_logfile=syslog
stderr_logfile=syslog