# docker-files/etc/supervisor/conf.d/supervisor-syslog-ng.conf
[program:syslogng]
command=/usr/sbin/syslog-ng --foreground --no-caps
stdout_logfile=/var/log/syslog-ng.log
redirect_stderr=true
priority=100